# services/veil-cli/veilkey-veilroot-egress-guard@.service
[Unit]
Description=VeilKey veilroot cgroup egress guard for profile %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=VEILKEY_VEILROOT_PROFILE=%i
ExecStart=/usr/local/bin/veilkey-veilroot-egress-guard %i watch
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
