# services/veil-cli/veilkey-veilroot-observe@.service
[Unit]
Description=VeilKey veilroot observer for profile %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=VEILKEY_VEILROOT_PROFILE=%i
Environment=VEILKEY_VEILROOT_WAIT_SECONDS=30
StandardOutput=append:/var/log/veilkey-proxy/veilroot-%i-observe.jsonl
StandardError=append:/var/log/veilkey-proxy/veilroot-%i-observe.log
ExecStart=/usr/local/bin/veilkey-veilroot-observe %i
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
